# sources/puppetcode_upgrade.timer
[Unit]
Description=Timer to Run Puppet Code Upgrade Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
