# sources/puppetcode_apply.timer
[Unit]
Description=Timer to Run Puppet Code Apply Every Hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
